# etc/cron.daily/backup-rsync
#!/bin/bash

# example to be put in /etc/cron.daily/
/opt/backup-rsync/backup.sh --daily --quiet