# etc/cron.weekly/tt-backup.sh
#!/bin/bash

# example to be put in /etc/cron.weekly/
/opt/backup-rsync/backup.sh --weekly --quiet